# roles/nginx/files/nginx.conf
description "Moxie Daemon container"
author "Paul R. Tagliamonte"
start on filesystem and started docker
stop on runlevel [!2345]
respawn

script

/usr/bin/docker start -a nginx || \
    /usr/bin/docker run \
        --name nginx \
        -p 0.0.0.0:80:80/tcp \
        -v /srv/docker/nginx/serve/:/serve \
        -v /srv/docker/nginx/sites-enabled/:/etc/nginx/sites-enabled/ \
        paultag/nginx \
        /usr/sbin/nginx \
            -c /etc/nginx/nginx.conf \
            -g "daemon off;"

end script

pre-stop exec /usr/bin/docker stop -t 5 uwsgi
